# trained-finesse.service
[Unit]
Description=ydotool macro loop
After=network.target

[Service]
ExecStart=/bin/sh /home/mel/.local/bin/trained-finesse.sh

